# roles/node_exporter/files/systemd.service
[Unit]
Description=Node Exporter

[Service]
User=daemon
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=default.target
