# roles/bird_exporter/files/bird_exporter.service
[Unit]
Description=Bird Exporter

[Service]
User=root
ExecStart=/usr/local/bin/bird_exporter -bird.ipv4=false

[Install]
WantedBy=default.target
